# backend/myah/cli/templates/hermes-gateway.service.in
[Unit]
Description=Hermes Agent gateway (Myah OSS)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=__HERMES_BIN__ gateway run
Restart=on-failure
RestartSec=5s
Environment=HERMES_HOME=__HERMES_HOME__

[Install]
WantedBy=default.target
